# supervisor_virasana_celery.conf
[program:virasana_celery]
directory=//home/ivan/ajna/virasana
command=//home/ivan/ajna/virasana/virasana-venv/bin/celery -A virasana.workers.tasks.celery worker --loglevel=info -n virasana@%%h
autostart=true
autorestart=true
user=ivan
stderr_logfile=/var/log/virasana/celery_error.log
stdout_logfile=/var/log/virasana/celery.log
startsecs=10