# supervisor_virasana.conf
[program:virasana]
directory=///home/ivan/ajna/virasana
command=///home/ivan/ajna/virasana/virasana-venv/bin/gunicorn wsgi_production:application -b localhost:5001
autostart=true
autorestart=true
user=ivan
stderr_logfile=/var/log/virasana/error.log
stdout_logfile=/var/log/virasana/access.log
